# systemd/opentofu-runner.service
[Unit]
PartOf=opentofu-runner.target
[Install]
WantedBy=opentofu-runner.target
[Service]
ExecStart=/usr/bin/ruby -e 'sleep'
Slice=manageiq.slice
